# etc/supervisor/conf.d/composer.conf

[program:openbox]
priority=50
directory=/home/balloon
; --sm-disable
;   Do not connect to the session manager.
command=/usr/bin/openbox --sm-disable --config-file /etc/X11/openbox/rc.xml --startup /composer/openbox/autostart.sh
user=balloon
autostart=true
autorestart=true
startsecs=2
stopsignal=QUIT
redirect_stderr=true

;[program:x11vnc]
;priority=20
;directory=/home/balloon
;command=x11vnc -display :1 -xkb -forever -xrandr -scale 2/3
;-norepeat              Option -norepeat disables X server key auto repeat when
;-repeat                VNC clients are connected and VNC keyboard input is
;                       not idle for more than 5 minutes.  This works around a
;                       repeating keystrokes bug (triggered by long processing
;                       delays between key down and key up client events:
;                       either from large screen changes or high latency).
;                       Default: -norepeat
; -nowirecopyrect
;
; command=x11vnc -display :0 -xkb -noxinerama -shared -xrandr -nolookup  -nopw -forever -unixsockonly /tmp/.unixx11vnc
;command=/composer/x11vnc-composer.sh
;user=balloon
;autostart=true
;autorestart=true
;stopsignal=QUIT
;startsecs=2
;stdout_logfile=/home/balloon/.logs/x11vnc.log
;redirect_stderr=true

[program:tigervnc]
priority=1
# priority=20
# directory=/home/balloon
# command=Xvnc :0 -geometry 3840x2160 -depth 24 -rfbunixpath /tmp/.x11vnc -pn -SecurityTypes None
# command=Xvnc :0 -geometry 3840x2160 -depth 24 -rfbport 5900 -pn -SecurityTypes None
# command=Xvnc :0 -geometry 3840x2160 -depth 24 -rfbport 5900 -pn -rfbauth /composer/run/.vnc/passwd
# -SendPrimary=0 disable the Linux X primary selection gets sent to the clipboard.
command=Xvnc :0 -geometry 3840x2160 -SendPrimary=0 -depth 24 -rfbunixpath /tmp/.x11vnc -pn -rfbauth /composer/run/.vnc/passwd
user=balloon
autostart=true
autorestart=true
stopsignal=QUIT
startsecs=2
stdout_logfile=/var/log/desktop/Xvnc.log
redirect_stderr=true

[program:novnc]
priority=25
# priority=25
directory=/composer
# command=ws-tcp-bridge --method=ws2tcp --lport 6081 --rhost=127.0.0.1:5900
# command=/usr/bin/python  /composer/noVNC/utils/websockify/websockify.py --unix-target /tmp/.x11vnc 6081
# command =/usr/bin/ws-tcp-bridge --method=ws2tcp --lport 6081 --rhost=unix:/tmp/.x11vnc
command =/composer/node/ws-tcp-bridge/ws-tcp-bridge --method=ws2tcp --lport 6081 --rhost=unix:/tmp/.x11vnc
environment=CONTAINER_IP_ADDR="%(ENV_CONTAINER_IP_ADDR)s"
;command=/composer/websockify-composer.sh
user=balloon
autostart=true
autorestart=true
startsecs=2
stopsignal=QUIT
redirect_stderr=true

[program:spawner-service]
priority=2
command=node /composer/node/spawner-service/spawner-service.js
user=balloon
autostart=true
autorestart=true
stopsignal=QUIT
redirect_stderr=true
;environment=NODE_EXTRA_CA_CERTS=/etc/ssl/certs/xxxx-root.pem,HOME=/home/balloon,LOGNAME=balloon
environment=CONTAINER_IP_ADDR="%(ENV_CONTAINER_IP_ADDR)s",KUBERNETES_SERVICE_HOST="%(ENV_KUBERNETES_SERVICE_HOST)s",SPAWNER_SERVICE_TCP_PORT=29786,BROADCAST_SERVICE_TCP_PORT=29784,HOME=/home/balloon,LOGNAME=balloon

[program:printer-service]
priority=90
directory=/home/balloon
command=node /composer/node/printer-service/printer-service.js
environment=CONTAINER_IP_ADDR="%(ENV_CONTAINER_IP_ADDR)s",BROADCAST_SERVICE_TCP_PORT=29784
user=balloon
autostart=true
autorestart=true
stopsignal=QUIT
redirect_stderr=true

[program:broadcast-service]
priority=5
directory=/home/balloon
command=node /composer/node/broadcast-service/broadcast-service.js
environment=CONTAINER_IP_ADDR="%(ENV_CONTAINER_IP_ADDR)s",BROADCAST_SERVICE_TCP_PORT=29784
user=balloon
autostart=true
autorestart=true
stopsignal=QUIT
redirect_stderr=true

[program:file-service]
priority=90
directory=/home/balloon
command=node /composer/node/file-service/file-service.js
environment=CONTAINER_IP_ADDR="%(ENV_CONTAINER_IP_ADDR)s",FILE_SERVICE_TCP_PORT=29783
user=balloon
autostart=true
autorestart=true
stopsignal=QUIT
redirect_stderr=true

[program:dbus-daemon-system]
priority=90
command=sudo /usr/bin/dbus-daemon --config-file=/usr/share/dbus-1/system.conf --print-address --fork --nosyslog
user=balloon
autostart=false
autorestart=true
stopsignal=QUIT
redirect_stderr=true

[program:xterm.js]
priority=90
directory=/composer/node/xterm.js
command=node /composer/node/xterm.js/main.js
user=balloon
autostart=true
autorestart=true
stopsignal=QUIT
environment=CONTAINER_IP_ADDR="%(ENV_CONTAINER_IP_ADDR)s", XTERM_TCP_PORT=29781, SHELL="/bin/bash", TERM="xterm-256color", LS_COLORS="rs=0:di=01;34:ln=01;36:mh=00:pi=40;33:so=01;35:do=01;35:bd=40;33;01:cd=40;33;01:or=40;31;01:su=37;41:sg=30;43:ca=30;41:tw=30;42:ow=34;42:st=37;44:ex=01;32:*.tar=01;31:*.tgz=01;31:*.arj=01;31:*.taz=01;31:*.lzh=01;31:*.lzma=01;31:*.tlz=01;31:*.txz=01;31:*.zip=01;31:*.z=01;31:*.Z=01;31:*.dz=01;31:*.gz=01;31:*.lz=01;31:*.xz=01;31:*.bz2=01;31:*.bz=01;31:*.tbz=01;31:*.tbz2=01;31:*.tz=01;31:*.deb=01;31:*.rpm=01;31:*.jar=01;31:*.war=01;31:*.ear=01;31:*.sar=01;31:*.rar=01;31:*.ace=01;31:*.zoo=01;31:*.cpio=01;31:*.7z=01;31:*.rz=01;31:*.jpg=01;35:*.jpeg=01;35:*.gif=01;35:*.bmp=01;35:*.pbm=01;35:*.pgm=01;35:*.ppm=01;35:*.tga=01;35:*.xbm=01;35:*.xpm=01;35:*.tif=01;35:*.tiff=01;35:*.png=01;35:*.svg=01;35:*.svgz=01;35:*.mng=01;35:*.pcx=01;35:*.mov=01;35:*.mpg=01;35:*.mpeg=01;35:*.m2v=01;35:*.mkv=01;35:*.webm=01;35:*.ogm=01;35:*.mp4=01;35:*.m4v=01;35:*.mp4v=01;35:*.vob=01;35:*.qt=01;35:*.nuv=01;35:*.wmv=01;35:*.asf=01;35:*.rm=01;35:*.rmvb=01;35:*.flc=01;35:*.avi=01;35:*.fli=01;35:*.flv=01;35:*.gl=01;35:*.dl=01;35:*.xcf=01;35:*.xwd=01;35:*.yuv=01;35:*.cgm=01;35:*.emf=01;35:*.axv=01;35:*.anx=01;35:*.ogv=01;35:*.ogx=01;35:*.aac=00;36:*.au=00;36:*.flac=00;36:*.mid=00;36:*.midi=00;36:*.mka=00;36:*.mp3=00;36:*.mpc=00;36:*.ogg=00;36:*.ra=00;36:*.wav=00;36:*.axa=00;36:*.oga=00;36:*.spx=00;36:*.xspf=00;36:"
redirect_stderr=true

[program:xsettingsd]
priority=90
directory=/home/balloon
command=/usr/bin/xsettingsd -c /home/balloon/.xsettings
user=balloon 
autostart=true
autorestart=true
stopsignal=QUIT
startsecs=10
redirect_stderr=true
