# etc/supervisor/conf.d/pulseaudio.conf
[program:pulseaudio]
priority=90
command=/usr/bin/pulseaudio
autostart=false
autorestart=true
stopsignal=QUIT
redirect_stderr=true

